# agent/systemd/infostealer.service
[Unit]
Description=Infostealer Simulation

[Service]
ExecStart=/opt/infostealer/infostealer.sh

[Install]
WantedBy=multi-user.target